# ansible/project/files/brook-template copy.service
[Unit]
Description=Brook from Aurora Admin Panel
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
Type=simple
User=root
DynamicUser=true
ExecStart=/usr/local/bin/brook

[Install]
WantedBy=multi-user.target
