# ansible/project/files/wstunnel-template.service
[Unit]
Description=Wstunnel from Aurora Admin Panel
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
Type=simple
User=root
DynamicUser=true
ExecStart=/usr/local/bin/wstunnel

[Install]
WantedBy=multi-user.target
